# deploy/firstboot/q-ui-firstboot.service
[Unit]
Description=3x-ui first-boot per-instance credential generation
Documentation=https://github.com/IntelligentQuantum/q-ui
# Run after the network and cloud-init are up, but BEFORE the panel starts, so
# the panel never serves the default admin/admin account.
After=network-online.target cloud-init.service
Wants=network-online.target
Before=q-ui.service
# Skip entirely once the sentinel exists (cheap guard; the script re-checks too).
ConditionPathExists=!/etc/q-ui/.firstboot-done

[Service]
Type=oneshot
RemainAfterExit=yes
# Inherit the same DB configuration the panel uses (sqlite default / postgres).
EnvironmentFile=-/etc/default/q-ui
EnvironmentFile=-/etc/conf.d/q-ui
EnvironmentFile=-/etc/sysconfig/q-ui
ExecStart=/usr/local/q-ui/q-ui-firstboot.sh

[Install]
WantedBy=multi-user.target
